# backend/crontab
7 1 * * * cd / && /usr/local/bin/python -m app.etl.extract > /var/log/extract
